# wifi-setup/systemd/wifi-failover.service
[Unit]
Description=WiFi failover monitor (real-Internet check, USB->native)
After=NetworkManager.service
Wants=NetworkManager.service
ConditionPathExists=/opt/wifi-setup/bin/wifi-failover-monitor

[Service]
Type=oneshot
ExecStart=/opt/wifi-setup/bin/wifi-failover-monitor
Nice=10
# Keep transient failures from spamming; the timer will retry.
SuccessExitStatus=0 1
