# wifi-setup/systemd/wifi-failover.timer
[Unit]
Description=Run WiFi failover monitor every 20s
After=NetworkManager.service

[Timer]
# First run shortly after boot, then every 20 seconds.
OnBootSec=30
OnUnitActiveSec=20
AccuracySec=2
Unit=wifi-failover.service

[Install]
WantedBy=timers.target
