# aicli/podman/quadlets/aicli-memory-engine.container
[Unit]
Description=AICLI Memory Engine (Rust/SQLite)
After=aicli-pod.service

[Container]
Pod=aicli-pod
Image=localhost/aicli-memory-engine:latest
Volume=%h/.local/share/aicli:/data:rw
Volume=/tmp/aicli:/tmp/aicli:rw
Environment=DB_PATH=/data/memory/aicli.db
Environment=MEMORY_SOCK=/tmp/aicli/memory.sock
Environment=SNAPSHOT_DIR=/data/snapshots

[Service]
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
