# aicli/podman/quadlets/aicli-orchestrator.container
[Unit]
Description=AICLI Orchestrator
After=aicli-sentinel.service aicli-memory-engine.service

[Container]
Pod=aicli-pod
Image=localhost/aicli-orchestrator:latest
Volume=%h/.local/share/aicli:/data:rw
Volume=%h/.config/aicli:/config:ro
Volume=/tmp/aicli:/tmp/aicli:rw
Environment=SENTINEL_SOCK=/tmp/aicli/sentinel.sock
Environment=MEMORY_SOCK=/tmp/aicli/memory.sock
Environment=ORCHESTRATOR_SOCK=/tmp/aicli/orchestrator.sock
Environment=DATA_DIR=/data

[Service]
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
