# aicli/podman/quadlets/aicli-compressor.container
[Unit]
Description=AICLI Context Compressor
After=aicli-pod.service

[Container]
Pod=aicli-pod
Image=localhost/aicli-compressor:latest
Volume=/tmp/aicli:/tmp/aicli:rw
Environment=COMPRESSOR_SOCK=/tmp/aicli/compressor.sock

[Service]
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
